# supervisord.conf
[supervisord]
nodaemon=true
logfile=/build/data/supervisor.log
pidfile=/build/data/supervisor.pid

[unix_http_server]
file = /var/supervisor.sock

[supervisorctl]
serverurl=unix:///var/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:celery]
command=/usr/bin/celery worker -D -E -A b2share.celery -l INFO --workdir=/eudat/b2share
environment=PYTHONPATH="/build/b2share:/build/b2share/demo"
autostart=false
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:celery_beat]
command=/usr/bin/celery beat -A b2share.celery --pidfile= --workdir=/eudat/b2share --loglevel="DEBUG"
environment=PYTHONPATH="/build/b2share:/build/b2share/demo"
autostart=false
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:wdb]
command=/usr/bin/wdb.server.py
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
